# systemd/scalibur-dashboard.service
[Unit]
Description=Scalibur Dashboard
After=network.target

[Service]
Type=simple
User=g
WorkingDirectory=/home/g/scalibur
ExecStart=/home/g/.local/bin/uv run python dashboard.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
