# systemd/scalibur-scanner.service
[Unit]
Description=Scalibur BLE Scale Scanner
After=bluetooth.target

[Service]
Type=simple
User=g
WorkingDirectory=/home/g/scalibur
ExecStart=/home/g/.local/bin/uv run python scanner.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
